# sonicroom.service
[Unit]
Description=SonicRoom audio conferencing server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/sonicroom
ExecStart=/usr/bin/pnpm start
Restart=on-failure
RestartSec=3

# Set your VPS public IP here
Environment=PORT=3100
Environment=ANNOUNCED_IP=217.154.98.41
Environment=ANNOUNCED_IP6=2001:ba0:208:5400::1
Environment=NODE_ENV=production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/sonicroom
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
